# modules/conf/dnspr_agency.conf
[program:dnspr_agency]
directory=/vagrant/modules/agencies

command=python server.py "http://mstrong.info/api/views/dnspr.json" "http://localhost:90/api/useragencyfields" 5
autostart=true
autorestart=true

